# config/throttle-me-daemon.service
[Unit]
Description=throttle-me hotspot detection daemon
Documentation=https://github.com/wtyler/throttle-me
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/throttle-me-daemon
Restart=always
RestartSec=10
Environment="DISPLAY=:0"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=throttle-me-daemon

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
